# deployment/services/oraculo-news.service
[Unit]
Description=Oráculo Bot - Servicio de Noticias
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/oraculo_bot
ExecStart=/path/to/your/oraculo_bot/venv/bin/python run_news_service.py
Restart=always
RestartSec=10

# Variables de entorno
Environment="SERVICE_MODE=news"
EnvironmentFile=-/path/to/your/oraculo_bot/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oraculo-news

# Límites de recursos
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target 